# config/cron/production.crontab
############# unitedstates tasks ###########

# legislators and committees from unitedstates
51 3 * * *  ~/sansad/config/cron/scripts/task.sh legislators 

# bills from THOMAS, once an hour (THOMAS scraper uses --fast mode)
39 * * * * ~/sansad/config/cron/scripts/task.sh bills limit=10


# leave the newline here
